# php-laravel/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autorestart=true
